# osconfig-agent-test.conf
description "OsConfig agent test"
start on stopped rc RUNLEVEL=[2345]

respawn
exec /usr/bin/osconfig_agent_test
